# kiri/packaging/macos/dev.kiri.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>             <string>dev.kiri</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/kiri</string>
    <string>serve</string>
  </array>
  <key>UserName</key>          <string>_kiri</string>
  <key>RunAtLoad</key>         <true/>
  <key>KeepAlive</key>         <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>KIRI_UPSTREAM_KEY_FILE</key> <string>/var/lib/kiri/upstream.key</string>
    <key>KIRI_CONFIG</key>            <string>/var/lib/kiri/config.yaml</string>
    <key>WORKSPACE</key>              <string>/var/lib/kiri/workspace</string>
  </dict>
  <key>StandardOutPath</key>   <string>/var/log/kiri/kiri.log</string>
  <key>StandardErrorPath</key> <string>/var/log/kiri/kiri.log</string>
</dict>
</plist>
